# scripts/boilerPlate.sh
#!/bin/sh

### BEGIN INIT INFO
# Provides:          serviceName
# Required-Start:    $network $remote_fs
# Required-Stop:     $network $remote_fs
# Should-Start:      avahi
# Should-Stop:       avahi
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Servince name
# Description:       Service description
### END INIT INFO
